# ecseq.conf
description "Starting Lene's ECSeq website"
author "Endre Bakken Stovner"

start on runlevel [2345]

exec cd /home/endrebak/code/lene_side/rna; /home/endrebak/local/bin/lein run